# crontab
* * * * * /bin/bash -l -c "cd /home/thermouser/thermofidz-client; bash worker.sh"
